# system/usr/share/ratflow/app_configs/90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add",
RUN+="/bin/bash -c '/bin/chmod o+w /sys/class/backlight/*/brightness'"

